# deploy/systemd/handpdf-api.service.example
[Unit]
Description=AI Handwritten PDF Recreator API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/AI-handwritten-PDF-Recreator/apps/api
EnvironmentFile=/etc/handpdf.env
ExecStart=/home/ubuntu/AI-handwritten-PDF-Recreator/apps/api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
